# cookbooks/nfsn/templates/default/nfsn-dns.service.erb
[Unit]
Description=Update NearlyFreeSpeech.NET DNS for colbyolson.com

[Service]
Type=oneshot
IgnoreSIGPIPE=false
ExecStart=/home/colby/nfsn.api.sh
